# Linux-64/supervisor/supervisord.d/kt_proxy.ini
[program:ktproxy]
stopasgroup=true
user=root
directory=/root/kt_proxy
command=/root/kt_proxy/ktproxy
redirect_stderr=true
stdout_logfile=/root/kt_proxy/supervisor/out.log
stderr_logfile=/root/kt_proxy/supervisor/out.err